# user/resume-local.target
[Unit]
Description=Resume target
AllowIsolate=true
After=root-resume.service

[Install]
WantedBy=suspend.target