# user/suspend-local.target
[Unit]
Description=Suspend target
AllowIsolate=true
Before=root-resume.service

[Install]
WantedBy=sleep.target