# user/udiskie.service
[Unit]
Description=UDiskie

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=session.target