# user/wm.target
[Unit]
Description=Window manager
Requires=dbus.socket
AllowIsolate=true

[Install]
WantedBy=default.target