# user/xbindkeys.service
[Unit]
Description=XBindKeys

[Service]
Type=simple
ExecStart=/usr/bin/xbindkeys

[Install]
WantedBy=session.target